# docker/root/etc/services.d/url-shortener/run
#!/usr/bin/with-contenv bash

# Start url-shortener
cd /app

export DANGEROUSLY_DISABLE_HOST_CHECK=true

exec \
    yarn start
